# server/powrush-server.service
[Unit]
Description=Powrush-MMO Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/powrush-server
WorkingDirectory=/opt/powrush
User=powrush
Group=powrush

# === Strong Hardening ===
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectProc=invisible

# Capability Bounding (recommended)
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW

# Filesystem
ReadWritePaths=/opt/powrush/data /opt/powrush/logs /tmp

# Network
PrivateNetwork=false

# seccomp
SystemCallFilter=@system-service
SystemCallFilter=~@debug @mount @cpu-emulation @obsolete @raw-io @reboot @swap @privileged

LimitNOFILE=65536
LimitNPROC=512

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
