# server/powrush-server.service
[Unit]
Description=Powrush-MMO Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/powrush-server
WorkingDirectory=/opt/powrush
User=powrush
Group=powrush

# === Strong Hardening ===
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectProc=invisible

# Filesystem
ReadWritePaths=/opt/powrush/data /opt/powrush/logs /tmp

# Network Namespace (strong isolation)
PrivateNetwork=false          # Set to true if using a reverse proxy or specific network setup
# NetworkNamespacePath=/run/netns/powrush   # Optional: use a dedicated network namespace

# seccomp
SystemCallFilter=@system-service
SystemCallFilter=~@debug @mount @cpu-emulation @obsolete @raw-io @reboot @swap @privileged

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
